# deploy/outpost-api-docker.service
[Unit]
Description=Outpost API (Docker Container)
After=docker.service
Requires=docker.service

[Service]
Type=simple
EnvironmentFile=/etc/outpost-api-docker.env
# Optional container envs (JWT secret, DB, etc.) are read by Docker from this file
# Keep systemd vars (IMAGE, API_PORT) separate from container envs
ExecStartPre=-/usr/bin/docker stop outpost-api
ExecStartPre=-/usr/bin/docker rm outpost-api
ExecStartPre=-/usr/bin/docker pull ${IMAGE}
ExecStart=/usr/bin/docker run --name outpost-api \
  -p ${API_PORT}:8080 \
  --env-file /etc/outpost-api.env \
  ${IMAGE}
ExecStop=/usr/bin/docker stop outpost-api
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target

